# scripts/greenlight-auto-deployer.service
[Unit]
Description=Greenlight Auto Deployer

[Service]
ExecStart=/bin/bash /usr/local/bin/greenlight-deploy
